# com.aura.cognitive.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.aura.cognitive.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/bryan/.gemini/antigravity/scratch/autonomy_engine/main_daemon.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/tmp/aura_daemon.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/aura_daemon.stderr.log</string>
    <key>WorkingDirectory</key>
    <string>/Users/bryan/.gemini/antigravity/scratch/autonomy_engine</string>
</dict>
</plist>
